# infra/nazar-api.service
[Unit]
Description=Nazar App FastAPI Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=tucibeyin
Group=tucibeyin
WorkingDirectory=/var/www/nazar_app
EnvironmentFile=/var/www/nazar_app/.env
ExecStart=/var/www/nazar_app/venv/bin/gunicorn -c /var/www/nazar_app/gunicorn.conf.py main:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nazar

# Güvenlik kısıtlamaları
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/nazar_app

[Install]
WantedBy=multi-user.target
